# wthr.timer
[Unit]
Description=timer for wthr.service for every hour

[Timer]
OnBootSec=60min
Unit=wthr.service

[Install]
WantedBy=timers.target
